# com.aitrading.autonomous.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.aitrading.autonomous</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/michaelmote/Desktop/ai-trading-system-complete/core/pacific_time_schedule.py</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    
    <key>WorkingDirectory</key>
    <string>/Users/michaelmote/Desktop/ai-trading-system-complete</string>
    
    <key>StandardOutPath</key>
    <string>/Users/michaelmote/Desktop/ai-trading-system-complete/logs/autonomous_system.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/michaelmote/Desktop/ai-trading-system-complete/logs/autonomous_system_error.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
        <key>PYTHONPATH</key>
        <string>/Users/michaelmote/Desktop/ai-trading-system-complete</string>
    </dict>
    
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>4</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    
    <key>UserName</key>
    <string>michaelmote</string>
    
    <key>GroupName</key>
    <string>staff</string>
</dict>
</plist>